# vps_server/systemd/fetch-iabp-data.service
[Unit]
Description=Fetch IABP buoy time-series data
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/arct-collector
ExecStart=/opt/arct-collector/.venv/bin/python3 /opt/arct-collector/scripts/fetch_iabp_data.py
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
