# vps_server/systemd/fetch-svalmiz-data.service
[Unit]
Description=Fetch SvalMIZ 2026 thermistor string buoy data from Thredds
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/arct-collector
ExecStart=/opt/arct-collector/.venv/bin/python3 /opt/arct-collector/scripts/fetch_svalmiz_data.py
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1
